# prod/etc/cron.daily/mongoalldump
#!/bin/sh
mongodump -o /work/PERSISTENT_DATA_DIR/mongodumps/`date "+%Y-%m-%d"`